# dist/debian/postinst
#!/bin/sh
set -e

case "$1" in
    configure)
        if ! getent passwd mrpump > /dev/null; then
            adduser --system --quiet --home /var/lib/mrpump --no-create-home --shell /bin/false --group mrpump
        fi
        mkdir -p /var/lib/mrpump
        chown mrpump:mrpump /var/lib/mrpump
        
        # Reload systemd and enable service
        if [ -d /run/systemd/system ]; then
            systemctl daemon-reload
            systemctl enable mrpump.service
        fi
    ;;
esac

exit 0
